# roles/hypervisor/templates/start-vm-domains.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/start-vms-that-are-shut-off.sh

[Install]
WantedBy=default.target
